# puppet/files/etc_supervisor/1-socorro-processor.conf
[program:socorro-processor]
user = socorro
command = /data/socorro/socorro/processor/processor_app.py --admin.conf=/etc/socorro/processor.ini
environment = PYTHONPATH="/data/socorro/application:/data/socorro/thirdparty"
autorestart = true
stderr_logfile = /var/log/socorro/processor-stderr.log
stdout_logfile = /var/log/socorro/processor-stdout.log
